# com.scorekeeper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.scorekeeper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mthirumalai/dev/scorekeeper/venv/bin/python3</string>
        <string>-m</string>
        <string>scorekeeper.main</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mthirumalai/dev/scorekeeper</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>8</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/mthirumalai</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mthirumalai/Library/Logs/scorekeeper.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mthirumalai/Library/Logs/scorekeeper.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
